# ygege/rootfs/etc/services.d/ygege/run
#!/usr/bin/with-contenv bashio

# Log starting
bashio::log.info "Starting Ygégé..."

# Ygégé doesn't have a complex config file by default, 
# it often relies on defaults or environment variables.
# If you need to pass specific ENV vars, you can export them here
# based on the options in config.yaml.

# Example: setting log level from config
export RUST_LOG=$(bashio::config 'log_level')

# Run the binary
# We exec so it takes over the PID
exec /usr/local/bin/ygege