# ygege/rootfs/etc/services.d/ygege/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Ygégé..."
exec /app/ygege
